# gentoo/etc/init.d/plexconnect
#!/sbin/runscript


depend() {
	need plex-media-server
	after plex-media-server
}

start() {
	ebegin "Starting PlexConnect"
	start-stop-daemon -S -m -p $PLEXCONNECT_PID -1 $PLEXCONNECT_LOG -2 $PLEXCONNECT_ERROR_LOG --quiet -u root -N -5 -b -d $PLEXCONNECT_PATH --exec $PLEXCONNECT_PATH/PlexConnect.py
	eend $?
}

stop() {
	ebegin "Stopping PlexConnect"
	kill -- -`cat $PLEXCONNECT_PID`
	rm -f $PLEXCONNECT_PID
	eend $?
}

